# reading-journal.service
[Unit]
Description=Reading Journal Application
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/reading-journal
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/reading-journal

[Install]
WantedBy=multi-user.target